# 1.0/config/supervisord.conf
[supervisord]
nodaemon=true

[program:apiserver]
autorestart =  true
autostart   =  true
command     =  bash -c "/usr/sbin/riak console"